# client/led-controller.service
[Unit]
Description=LED Controller Service
After=network.target

[Service]
Type=simple
User=brent
WorkingDirectory=/home/brent/dragonball-light/client
ExecStart=/home/brent/dragonball-light/client/venv/bin/python /home/brent/dragonball-light/client/led_controller.py --server https://dragonball.mittn.ca
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target